# homegear/base/image-data/etc/services.d/homegear/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid homegear /usr/bin/homegear -c /config/homegear