# homegear/base/image-data/etc/services.d/homegear-influxdb/run
#!/usr/bin/with-contenv bash

/usr/bin/homegear-influxdb -c /config/homegear
