# homegear/base/image-data/etc/services.d/logs/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid homegear /scripts/logs-tailer.sh