# homegear/base/image-data/etc/services.d/homegear-management/run
#!/usr/bin/with-contenv bash

/usr/bin/homegear-management -c /config/homegear