# homegear/base/image-data/etc/services.d/homegear-webssh/run
#!/usr/bin/with-contenv bash

/usr/bin/homegear-webssh -c /config/homegear
